# deploy/cleanup_expired_shares.timer.example
[Unit]
Description=Run cleanup_expired_shares daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
